# .docker/supervisord.ini
[supervisord]
nodaemon=true

[program:idrac-fan-control-daemon]
command=php /app/idrac-fan-control fans:daemon
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
